# taigabot.service
[Unit]
Description=TaigaBot Lite — Telegram bot with daily AI summaries
After=network.target

[Service]
User=botuser
WorkingDirectory=/home/botuser/bots/taigabot
ExecStart=/home/botuser/bots/taigabot/venv/bin/python bot.py
EnvironmentFile=/home/botuser/bots/taigabot/.env
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
